# config/supervisor.conf
[program:md-resume]
command=/usr/local/bin/gunicorn wsgi -c config/gunicorn.py
directory=/var/www/md-resume
autostart=true
autorestart=true
