# config/supervisor.conf
[program:md_resume]
command=/usr/local/bin/gunicorn wsgi -c config/gunicorn.py
directory=/var/www/md_resume
autostart=true
autorestart=true
